# deploy/darkroom-backup.timer
[Unit]
Description=Run darkroom catalog backup nightly

[Timer]
# 01:30, not later: the NAS powers down at 02:00 and wakes at 08:00, so the
# push window closes shortly after this fires.
OnCalendar=*-*-* 01:30:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
